# supervisor/amazon-sqs-worker.conf
[program:amazon-sqs-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/artisan queue:work sqs --sleep=3 --tries=3
autostart=true
autorestart=true
;user=root
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/storage/logs/amazon-sqs-worker.log
